# cassandra-repair-scheduler/src/deb/control/postrm
#!/bin/sh
if [ $(cat /proc/1/comm) = "init" ]
then
    update-rc.d cassandra-repair-scheduler remove
fi
